# etc/supervisord.d/check_location.ini
[program:check_location]
command = /usr/bin/perl /usr/local/sbin/project/check_location.pl
directory = /usr/local/sbin/project/
process_name = %(program_name)s-%(process_num)s
numprocs = 1
autorestart = true
startsecs = 3
user = root
log_stderr = true
logfile = /var/log/project/check_location_std.log
logfile_maxbytes = 200MB
