# etc/supervisord.d/reg_fail.ini
[program:reg_fail]
command = /usr/bin/perl /usr/local/sbin/project/reg_fail.pl
directory = /usr/local/sbin/project/
process_name = %(program_name)s-%(process_num)s
numprocs = 1
autorestart = true
startsecs = 3
user = root
log_stderr = true
logfile = /var/log/project/reg_fail_std.log
logfile_maxbytes = 200MB
