# etc/supervisord.d/register_repair.ini
[program:register_repair]
command = /usr/bin/perl /usr/local/sbin/project/register_repair.pl
directory = /usr/local/sbin/project/
process_name = %(program_name)s-%(process_num)s
numprocs = 1
autorestart = true
startsecs = 3
user = root
log_stderr = true
logfile = /var/log/project/register_repair_std.log
logfile_maxbytes = 200MB
